# bash/scraps/nodejs.conf
description "node.js server"
author      "kvz"

start on startup
stop on shutdown

script
    # We found $HOME is needed. Without it, we ran into problems
    export HOME="/root"

    exec /usr/local/bin/node /where/yourprogram.js 2>&1 >> /var/log/node.log
end script